# XyvorraOS/airootfs/etc/systemd/system/xyvorra-reboot-check.service
[Unit]
Description=XyvorraOS Reboot Configuration Check
After=multi-user.target
ConditionPathExists=/etc/xyvorra-liveusb-configured

[Service]
Type=oneshot
ExecStart=/usr/local/bin/xyvorra-reboot-check
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
